# .config/autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=copyq --start-server
Terminal=false
