# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst
Comment=Dunst is a notification server
Exec=dunst
Terminal=false
